# disabled/utils/utils-installer/assets/assets/os-common/harbor-update.service
[Unit]
Description=Harbor System Config Update Service

[Service]
ExecStart=/var/usrlocal/bin/harbor-update.sh
Type=oneshot
