# disabled/utils/utils-installer/assets/assets/os-common/harbor-update.path
[Unit]
Description=Harbor System Config Update Service

[Path]
PathModified=/tmp/harbor/harbor-update

[Install]
WantedBy=multi-user.target
