# disabled/utils/utils-installer/assets/assets/os-common/harbor-mounter.path
[Unit]
Description=HarborOS Mounter

[Path]
PathChanged=/tmp/harbor/harbor-mount

[Install]
WantedBy=multi-user.target
