# docker/v3/keystone/undercloud-cockpit/assets/usr/lib/systemd/system/dummy-service.service
[Unit]
Description=Dummy Service

[Service]
Type=oneshot
RemainAfterExit=True
ExecStart=/bin/true

[Install]
WantedBy=multi-user.target
